# etc/update-motd.d/98-reboot-required
#!/bin/sh

running=$(uname -r)
latest=$(find /boot/vmlinuz-* -type f 2>/dev/null | sort -V | tail -1 | sed 's|/boot/vmlinuz-||')

if [ -n "$latest" ] && [ "$running" != "$latest" ]; then
    printf " \033[91mSystem restart required (kernel: %s -> %s)\033[0m\n\n" "$running" "$latest"
else
    printf " No restart required (kernel: %s)\n\n" "$running"
fi