# etc/systemd/system/motd-updates.timer
[Unit]
Description=Run motd-updates at boot and hourly

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h
Persistent=yes

[Install]
WantedBy=timers.target

