# etc/systemd/system/motd-updates.service
[Unit]
Description=Compute MOTD update counts

[Service]
Type=oneshot
ExecStart=/usr/local/bin/motd-updates.sh

